# Templates/Installer_Package_with_LD/simple-package/scripts/postinstall
#!/bin/bash
launchctl load "/Library/LaunchDaemons/com.simple.agent.plist"

exit 0
